# systemd/openclaw-agent.service
; OpenClaw Super Agent - Windows Service
; Install with: sc create OpenClawAgent binPath= "C:\path\to\openclaw-agent.exe" start= auto
; Delete with: sc delete OpenClawAgent

[Unit]
Description=OpenClaw Super Agent - Self-Upgrading Personal AI
Documentation=https://github.com/aimerdoux/super-agent-architecture
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=notify
# Notify systemd when ready (requires sd_notify support)
NotifyAccess=main
# Working directory
WorkingDirectory=%USERPROFILE%\.openclaw\workspace
# Executable path
ExecStart=C:\Users\play4\AppData\Roaming\npm\node_modules\openclaw\bin\openclaw.exe gateway start --port 18789
# Restart settings
Restart=always
RestartSec=10
# Environment variables
Environment=OPENCLAW_Workspace=%USERPROFILE%\.openclaw\workspace
Environment=SUPABASE_URL=https://your-project.supabase.co
Environment=SUPABASE_SERVICE_KEY=your-service-key
Environment=GITHUB_TOKEN=your-github-token
# Health check on startup
ExecStartPost=PowerShell -Command "& { Start-Sleep -Seconds 5; Invoke-RestMethod -Uri 'https://your-project.functions.supabase.co/agent-core' -Method Post -Body '{ \"action\": \"status\" }' -ContentType 'application/json' }"
# Graceful shutdown
ExecStop=PowerShell -Command "& { Invoke-RestMethod -Uri 'https://your-project.functions.supabase.co/agent-core' -Method Post -Body '{ \"action\": \"restart\" }' -ContentType 'application/json' }"
TimeoutStartSec=60
TimeoutStopSec=30

# Resource limits
MemoryMax=2G
CPUQuota=80%
FileDescriptorLimit=1024

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-agent

# Security (if running as service user)
; User=play4
; Group=users
; NoNewPrivileges=true
; ProtectSystem=strict
; ProtectHome=true
; PrivateTmp=true

[Install]
WantedBy=multi-user.target
